# install/macos/LaunchAgent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.driftwatch.riskagent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/driftwatch-agent</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/var/log/driftwatch-agent/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/driftwatch-agent/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>DRIFTWATCH_AGENT__LOG_LEVEL</key>
        <string>info</string>
    </dict>

    <!-- Required for CoreLocation access -->
    <key>NSLocationWhenInUseUsageDescription</key>
    <string>driftwatch-agent requires location access to compute workspace geo-anchor risk.</string>
</dict>
</plist>
